# iac/lxc/systemd/homelab-agent.service
[Unit]
Description=Foxhole homelab diagnostic agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=agent
Group=agent
WorkingDirectory=/opt/homelab-agent
EnvironmentFile=/etc/homelab-agent/foxhole.env
ExecStart=/opt/homelab-agent/venv/bin/uvicorn agent.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/homelab-agent /etc/homelab-agent
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
